# run_apps.conf
[program:run_apps]
command=/home/sanyash/TKP/run_apps.sh
autostart=true
autorestart=true
stdout_logfile=/var/log/app.log
redirect_stderr=true
